# scripts/invoke_news_please.service
[Unit]
Description=Invoke news-please shell script

[Service]
Type=oneshot
ExecStart=./invoke_news_please.sh
