# scripts/invoke_news_please.timer
[Unit]
Description=Run invoke_news_please.service every 1 hours

[Timer]
OnUnitActiveSec=1h
AccuracySec=1s
Unit=invoke_news_please.service

[Install]
WantedBy=timers.target
